# services/docker-cron/crontab
*/1 * * * * root python3 /usr/src/app/poll_metadata_and_obs_types.py >/proc/1/fd/1 2>/proc/1/fd/2
